# data/systemd/okerrclient.service
#
# 
# 

[Unit]
Description=okerr client daemon
After=syslog.target
After=network.target

[Service]
Type=simple
# User=<run as user>
# Group=<run as group group>
ExecStart=/usr/local/bin/okerrclient -d --fg
StandardOutput=null

# Give the script some time to startup
# TimeoutSec=300

[Install]
WantedBy=multi-user.target

